# docker/supervisor.conf.d/janitoo_tellstick.conf
[program:jnt_tellstick]
user=root
autostart=true
autorestart=true
command=jnt_tellstick -c /opt/janitoo/etc/jnt_tellstick.conf start
